# supervisord.conf
[program:robo-gong]
command=~/robo-gong/venv/bin/python ~/robo-gong/robo-gong.py
stderr_logfile = /var/log/supervisord/robo-gong.log
stdout_logfile = /var/log/supervisord/robo-gong.log
